# saha-layers/meta-tegra-saha/recipes-saha/ble-identity/saha-ble-identity/saha-ble-identity.service
[Unit]
Description=Initialize persistent Saha BLE identity
Documentation=file:///usr/share/doc/saha-bt-wifi-provision/GATT.md
Requires=tegra-bluetooth.service
After=tegra-bluetooth.service
Before=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/bin/saha-ble-identity-init
RemainAfterExit=yes

[Install]
WantedBy=bluetooth.target
