# saha-layers/meta-tegra-saha/recipes-saha/homeassistant-container/saha-homeassistant-container/homeassistant-container.service
[Unit]
Description=Home Assistant Container
Documentation=https://www.home-assistant.io/installation/linux#install-home-assistant-container
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/homeassistant-container
ExecStartPre=/usr/bin/saha-homeassistant-container wait-docker
ExecStart=/usr/bin/saha-homeassistant-container start
ExecStop=/usr/bin/saha-homeassistant-container stop
TimeoutStartSec=900
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
