# saha-layers/meta-tegra-saha/recipes-saha/bt-wifi-provision/saha-bt-wifi-provision/saha-bt-wifi-provision.service
[Unit]
Description=Saha Bluetooth WiFi provisioning GATT service
Documentation=file:///usr/share/doc/saha-bt-wifi-provision/GATT.md
After=bluetooth.service tegra-bluetooth.service NetworkManager.service NetworkManager-wait-online.service
Wants=tegra-bluetooth.service NetworkManager.service NetworkManager-wait-online.service
Requires=bluetooth.service
# Stop immediately if BlueZ disappears and follow every BlueZ stop/restart job.
BindsTo=bluetooth.service
PartOf=bluetooth.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/saha-bt-wifi-provision
Environment=PYTHONUNBUFFERED=1
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
StateDirectory=saha
StateDirectoryMode=0700
ReadWritePaths=/run /var/lib/NetworkManager
ExecStartPre=/usr/bin/saha-bt-wifi-provision-wait
ExecStart=/usr/bin/saha-bt-wifi-provision
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
