# saha-layers/meta-tegra-saha/recipes-connectivity/bluez/bluez5/bluetooth.service.d/saha-experimental.conf
[Unit]
Requires=saha-ble-identity.service
After=saha-ble-identity.service
# Pull the GATT application back in whenever bluetooth.service starts again.
Wants=saha-bt-wifi-provision.service
Before=saha-bt-wifi-provision.service
# Avoid an endless crash loop while still recovering transient bluetoothd faults.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --experimental
Restart=on-failure
RestartSec=3
# Match the directory installed by the BlueZ package and remove the mode warning.
ConfigurationDirectoryMode=0755
