# saha-layers/meta-tegra-saha/recipes-saha/docker-compose/saha-docker-compose/saha-docker-compose.service
[Unit]
Description=Saha Docker Compose application stack
Documentation=https://docs.docker.com/compose/
After=docker.service network-online.target systemd-timesyncd.service time-sync.target
Wants=network-online.target systemd-timesyncd.service time-sync.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/saha-docker-compose
ExecStartPre=/usr/bin/saha-docker-compose wait-docker
ExecStart=/usr/bin/saha-docker-compose start
ExecStop=/usr/bin/saha-docker-compose stop
TimeoutStartSec=900
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
